# distro/live-build/config/includes.chroot/usr/lib/systemd/system/kelp-gateway.service
[Unit]
Description=Kelp AI Gateway
Documentation=https://github.com/kelp-linux/kelp-linux
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/kelp-gateway
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
WatchdogSec=30

User=kelp
Group=kelp
RuntimeDirectory=kelp
StateDirectory=kelp
LogsDirectory=kelp
ConfigurationDirectory=kelp

# Hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=~user
SystemCallFilter=@system-service @sandbox
SystemCallErrorNumber=EPERM
MemoryDenyWriteExecute=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
